# deploy/io.goodkind.lmd.serve.test.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
    Template for an ISOLATED lmd-serve test LaunchAgent.

    This is NOT installed by hand. scripts/lmd-test-daemon.sh fills the
    double-brace placeholders below against the worktree's freshly built
    lmd-serve and bootstraps the result, so a branch can be exercised end to end
    without ever touching the production daemon.

    Isolation from production (io.goodkind.lmd.serve on :5400):
      - distinct launchd Label and a distinct pair of Mach services, so launchd
        never sees a name collision;
      - a distinct LMD_PORT and LMD_DATA_DIR;
      - LMD_HOST_SERVICE / LMD_CONTROL_SERVICE / LMD_LAUNCHD_LABEL tell the broker
        to register and expect the test names instead of the production ones.

    The battery thresholds are filled in by scripts/lmd-test-daemon.sh, which
    defaults them so the PowerMonitor is disabled for a test run. That keeps the
    hard-level admission halt from refusing requests on a laptop on battery; a
    run can override the defaults to exercise the real throttle.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LABEL}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{LMD_SERVE_PATH}}</string>
    </array>

    <key>MachServices</key>
    <dict>
        <key>{{CONTROL_SERVICE}}</key>
        <true/>
        <key>{{HOST_SERVICE}}</key>
        <true/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LMD_HOST</key>
        <string>localhost</string>
        <key>LMD_PORT</key>
        <string>{{LMD_PORT}}</string>
        <key>LMD_RESERVE_GB</key>
        <string>8</string>
        <key>LMD_SWIFTLM_BINARY</key>
        <string>{{LMD_SWIFTLM_BINARY}}</string>
        <key>LMD_CHAT_MAX_CONCURRENCY</key>
        <string>4</string>
        <key>LMD_EMBEDDING_MAX_CONCURRENCY</key>
        <string>4</string>

        <!-- Filled by lmd-test-daemon.sh. The default 0/1/2 disables the
             PowerMonitor (engage<=0) for the duration of a test run; override
             with LMD_TEST_BATTERY_* to exercise the real throttle. Values must
             satisfy hard < mild < resume. -->
        <key>LMD_BATTERY_THROTTLE_PCT</key>
        <string>{{LMD_BATTERY_THROTTLE_PCT}}</string>
        <key>LMD_BATTERY_MILD_PCT</key>
        <string>{{LMD_BATTERY_MILD_PCT}}</string>
        <key>LMD_BATTERY_RESUME_PCT</key>
        <string>{{LMD_BATTERY_RESUME_PCT}}</string>

        <key>LMD_DISABLE_XPC</key>
        <string>0</string>
        <key>LMD_IDLE_MINUTES</key>
        <string>15</string>
        <key>LMD_EMBEDDING_IDLE_MINUTES</key>
        <string>60</string>
        <key>LMD_DATA_DIR</key>
        <string>{{LMD_DATA_DIR}}</string>
        <key>LMD_SAMPLE_INTERVAL</key>
        <string>15</string>
        <key>LMD_PROMPT_CACHE_MAX_TOKENS</key>
        <string></string>
        <key>LMD_PROMPT_CACHE_ENABLED</key>
        <string>true</string>
        <key>LMD_MLX_CACHE_LIMIT_GB</key>
        <string>2</string>

        <!-- Isolation overrides: register and expect the test identities. -->
        <key>LMD_HOST_SERVICE</key>
        <string>{{HOST_SERVICE}}</string>
        <key>LMD_CONTROL_SERVICE</key>
        <string>{{CONTROL_SERVICE}}</string>
        <key>LMD_LAUNCHD_LABEL</key>
        <string>{{LABEL}}</string>

        <key>LMD_TRACE_DISABLE_MLX_SNAPSHOT</key>
        <string>0</string>

        <!-- OTLP export, filled by scripts/lmd-test-daemon.sh from the caller's
             OTEL_EXPORTER_OTLP_ENDPOINT. Empty (the default) disables export, so
             the broker and its spawned hosts stay on the in-process plane unless
             a run points them at a collector. The hosts inherit these from the
             broker's environment at spawn time. -->
        <key>OTEL_EXPORTER_OTLP_ENDPOINT</key>
        <string>{{OTEL_EXPORTER_OTLP_ENDPOINT}}</string>
        <key>OTEL_EXPORTER_OTLP_PROTOCOL</key>
        <string>{{OTEL_EXPORTER_OTLP_PROTOCOL}}</string>
        <key>OTEL_METRIC_EXPORT_INTERVAL</key>
        <string>{{OTEL_METRIC_EXPORT_INTERVAL}}</string>

        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>{{STDERR_LOG}}</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
