# deploy/io.goodkind.lmd.serve.test.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
    Template for an ISOLATED lmd-serve test LaunchAgent.

    This is NOT installed by hand. scripts/lmd-test-daemon.sh renders the
    {{PLACEHOLDERS}} below against the worktree's freshly built lmd-serve and
    bootstraps the result, so a branch can be exercised end to end without ever
    touching the production daemon.

    Isolation from production (io.goodkind.lmd.serve on :5400):
      - distinct launchd Label and a distinct pair of Mach services, so launchd
        never sees a name collision;
      - a distinct LMD_PORT and LMD_DATA_DIR;
      - LMD_HOST_SERVICE / LMD_CONTROL_SERVICE / LMD_LAUNCHD_LABEL tell the broker
        to register and expect the test names instead of the production ones.

    The battery throttle is disabled here (LMD_BATTERY_THROTTLE_PCT=0 turns the
    PowerMonitor off entirely), so the hard-level admission halt can never refuse
    requests on a laptop running on battery during a test.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LABEL}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{LMD_SERVE_PATH}}</string>
    </array>

    <key>MachServices</key>
    <dict>
        <key>{{CONTROL_SERVICE}}</key>
        <true/>
        <key>{{HOST_SERVICE}}</key>
        <true/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LMD_HOST</key>
        <string>localhost</string>
        <key>LMD_PORT</key>
        <string>{{LMD_PORT}}</string>
        <key>LMD_RESERVE_GB</key>
        <string>8</string>
        <key>LMD_SWIFTLM_BINARY</key>
        <string>{{LMD_SWIFTLM_BINARY}}</string>
        <key>LMD_CHAT_MAX_CONCURRENCY</key>
        <string>4</string>
        <key>LMD_EMBEDDING_MAX_CONCURRENCY</key>
        <string>4</string>

        <!-- Battery throttle OFF: engage<=0 disables PowerMonitor, so the hard
             admission halt never fires. mild/resume only need to satisfy the
             ordering validation (hard < mild < resume). -->
        <key>LMD_BATTERY_THROTTLE_PCT</key>
        <string>0</string>
        <key>LMD_BATTERY_MILD_PCT</key>
        <string>1</string>
        <key>LMD_BATTERY_RESUME_PCT</key>
        <string>2</string>

        <key>LMD_DISABLE_XPC</key>
        <string>0</string>
        <key>LMD_IDLE_MINUTES</key>
        <string>15</string>
        <key>LMD_EMBEDDING_IDLE_MINUTES</key>
        <string>60</string>
        <key>LMD_DATA_DIR</key>
        <string>{{LMD_DATA_DIR}}</string>
        <key>LMD_SAMPLE_INTERVAL</key>
        <string>15</string>
        <key>LMD_PROMPT_CACHE_MAX_TOKENS</key>
        <string></string>
        <key>LMD_PROMPT_CACHE_ENABLED</key>
        <string>true</string>
        <key>LMD_MLX_CACHE_LIMIT_GB</key>
        <string>2</string>

        <!-- Isolation overrides: register and expect the test identities. -->
        <key>LMD_HOST_SERVICE</key>
        <string>{{HOST_SERVICE}}</string>
        <key>LMD_CONTROL_SERVICE</key>
        <string>{{CONTROL_SERVICE}}</string>
        <key>LMD_LAUNCHD_LABEL</key>
        <string>{{LABEL}}</string>

        <key>LMD_TRACE_DISABLE_MLX_SNAPSHOT</key>
        <string>0</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>{{STDERR_LOG}}</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
