# deploy/io.goodkind.lmd.serve.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
    LaunchAgent for lmd-serve.
    Copy to ~/Library/LaunchAgents/ and load with:
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/io.goodkind.lmd.serve.plist
    Unload with:
      launchctl bootout gui/$(id -u)/io.goodkind.lmd.serve
    Restart to pick up a new binary:
      launchctl kickstart -k gui/$(id -u)/io.goodkind.lmd.serve

    lmd-serve is the single daemon that owns:
      - OpenAI-compatible broker on LMD_HOST:LMD_PORT
      - Sensor sampling (JSONL under LMD_DATA_DIR/memory.jsonl)
      - Fan control is disabled during the current moratorium; macOS owns fans
    This replaces the former swiftlmd + swiftmon pair.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.goodkind.lmd.serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{LMD_SERVE_PATH}}</string>
    </array>

    <!--
        Mach service for first-party Swift clients (lmd CLI, lmd-tui).
        launchd demand-launches the daemon when a client opens an
        XPCSession to this name. The daemon stays up under KeepAlive.
        Bumping this string is a breaking change; SwiftLMControl pins
        the same constant in BrokerProtocol.swift.
    -->
    <key>MachServices</key>
    <dict>
        <key>io.goodkind.lmd.control</key>
        <true/>
        <key>io.goodkind.lmd.host</key>
        <true/>
    </dict>

    <!--
        The broker fails fast at startup unless every key below is defined; see
        docs/configuration.md for the reference. There are no silent code
        defaults, so a missing key is a named startup error, not a guess.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>LMD_HOST</key>
        <string>localhost</string>
        <key>LMD_PORT</key>
        <string>5400</string>
        <key>LMD_RESERVE_GB</key>
        <string>20</string>
        <key>LMD_SWIFTLM_BINARY</key>
        <string>/Users/agoodkind/Sites/SwiftLM/.build/arm64-apple-macosx/release/SwiftLM</string>
        <key>LMD_CHAT_MAX_CONCURRENCY</key>
        <string>4</string>
        <key>LMD_EMBEDDING_MAX_CONCURRENCY</key>
        <string>1</string>
        <key>LMD_EMBED_BATCH_TOKEN_BUDGET</key>
        <string></string>
        <key>LMD_EMBED_BATCH_MAX_ROWS</key>
        <string>256</string>
        <key>LMD_EMBED_PRIORITY_MAX_INPUTS</key>
        <string>2</string>
        <key>LMD_EMBED_PRIORITY_MAX_TOKENS</key>
        <string>2048</string>
        <key>LMD_EMBED_PRIORITY_LANE</key>
        <string>true</string>
        <key>LMD_BATTERY_THROTTLE_PCT</key>
        <string>20</string>
        <key>LMD_BATTERY_MILD_PCT</key>
        <string>35</string>
        <key>LMD_BATTERY_RESUME_PCT</key>
        <string>80</string>
        <key>LMD_DISABLE_XPC</key>
        <string>0</string>
        <key>LMD_IDLE_MINUTES</key>
        <string>15</string>
        <key>LMD_EMBEDDING_IDLE_MINUTES</key>
        <string>60</string>
        <key>LMD_DATA_DIR</key>
        <string>/Users/agoodkind/Library/Application Support/io.goodkind.lmd</string>
        <key>LMD_SAMPLE_INTERVAL</key>
        <string>15</string>
        <key>LMD_PROMPT_CACHE_MAX_TOKENS</key>
        <string></string>
        <key>LMD_PROMPT_CACHE_ENABLED</key>
        <string>true</string>
        <key>LMD_MLX_CACHE_LIMIT_GB</key>
        <string></string>
        <key>LMD_TRACE_DISABLE_MLX_SNAPSHOT</key>
        <string>0</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!--
        Diagnostics live in unified logging under subsystem io.goodkind.lmd.
        Read them with the `log stream` and `log show` tools; README carries
        the exact flags. stdout routes to /dev/null so structured os.Logger
        output is the only normal path. stderr routes to a real file so Swift
        runtime crash banners (Fatal error, Precondition failed, native
        SIGSEGV) survive long enough to diagnose; the os.Logger pipeline does
        not capture stderr.
    -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/Users/agoodkind/Library/Logs/lmd-serve.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
